# com.utsavmehta.iplbot.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.utsavmehta.iplbot.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>--experimental-sqlite</string>
        <string>/Users/utsavmehta/PersonalProjects/gameday/dist/ipl-sync.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>UserName</key>
    <string>utsavmehta</string>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday/logs/ipl-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday/logs/ipl-sync-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
