# com.utsavmehta.f1bot.checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.utsavmehta.f1bot.checker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>--experimental-sqlite</string>
        <string>/Users/utsavmehta/PersonalProjects/gameday/dist/checker.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday</string>

    <!-- Run every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Run as this user (required for LaunchDaemon) -->
    <key>UserName</key>
    <string>utsavmehta</string>

    <key>StandardOutPath</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday/logs/checker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/utsavmehta/PersonalProjects/gameday/logs/checker-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
